# packaging/macos/scripts/postinstall
#!/bin/bash
set -uo pipefail

APP_SUPPORT="/Library/LiRAYS-Scada"
CONFIG="$APP_SUPPORT/settings.yaml"
DATA_DIR="$APP_SUPPORT/data_dir"
LOG_DIR="/Library/Logs/LiRAYS-Scada"
LAUNCHD_PLIST="/Library/LaunchDaemons/com.lirays.scada.plist"
PAYLOAD_PLIST_ROOT="/Library/LaunchDaemons/com.lirays.scada.root.plist"
BIN="/usr/local/bin/lirays-scada"
INSTALL_LOG="$LOG_DIR/install.log"

mkdir -p "$APP_SUPPORT" "$DATA_DIR" "$LOG_DIR"
chmod 755 "$APP_SUPPORT" "$DATA_DIR" "$LOG_DIR"
# Ensure ownership/permissions
chown -R root:wheel "$APP_SUPPORT" "$LOG_DIR"
[ -f "$LAUNCHD_PLIST" ] && chown root:wheel "$LAUNCHD_PLIST" && chmod 644 "$LAUNCHD_PLIST"
[ -f "$PAYLOAD_PLIST_ROOT" ] && chown root:wheel "$PAYLOAD_PLIST_ROOT" && chmod 644 "$PAYLOAD_PLIST_ROOT"
[ -f "$BIN" ] && chown root:wheel "$BIN" && chmod 755 "$BIN"
# Remove quarantine/xattrs that can block launchd
if command -v xattr >/dev/null 2>&1; then
  xattr -c "$BIN" "$LAUNCHD_PLIST" "$PAYLOAD_PLIST_ROOT" 2>/dev/null || true
fi

# Keep existing config if present; otherwise seed default from payload
if [ ! -f "$CONFIG" ] && [ -f "$CONFIG.default" ]; then
  mv "$CONFIG.default" "$CONFIG"
fi

# Log for troubleshooting
exec >>"$INSTALL_LOG" 2>&1
echo "=== postinstall $(date) ==="

# Ensure clean previous instance
launchctl bootout system/com.lirays.scada 2>/dev/null || true
launchctl disable system/com.lirays.scada 2>/dev/null || true

echo "Bootstrapping launchd service..."
# Prefer the root-owned plist if present (system domain requires root)
PLIST_TO_LOAD="$LAUNCHD_PLIST"
if [ -f "$PAYLOAD_PLIST_ROOT" ]; then
  mv -f "$PAYLOAD_PLIST_ROOT" "$LAUNCHD_PLIST"
  PLIST_TO_LOAD="$LAUNCHD_PLIST"
fi
set +e
BOOT_OUT=$(launchctl bootstrap system "$PLIST_TO_LOAD" 2>&1)
BOOT_RC=$?
set -e
echo "bootstrap rc=$BOOT_RC: $BOOT_OUT"

if [ $BOOT_RC -eq 0 ]; then
  launchctl enable system/com.lirays.scada || true
  launchctl kickstart -k system/com.lirays.scada || true
  echo "Service started."
else
  echo "WARNING: launchctl bootstrap failed (rc=$BOOT_RC). Trying legacy load -w..." >&2
  set +e
  LOAD_OUT=$(launchctl load -w "$PLIST_TO_LOAD" 2>&1)
  LOAD_RC=$?
  set -e
  echo "load rc=$LOAD_RC: $LOAD_OUT"
  if [ $LOAD_RC -eq 0 ]; then
    echo "Service loaded via launchctl load -w."
  else
    echo "ERROR: Service not started. Try manually: sudo launchctl bootstrap system \"$PLIST_TO_LOAD\"; sudo launchctl enable system/com.lirays.scada; sudo launchctl kickstart -k system/com.lirays.scada" >&2
  fi
fi

exit 0
