# debian/lirays-scada/DEBIAN/postinst
#!/bin/sh
# postinst script for lirays-scada.
#
# See: dh_installdeb(1).

set -e

# Summary of how this script can be called:
#        * <postinst> 'configure' <most-recently-configured-version>
#        * <old-postinst> 'abort-upgrade' <new version>
#        * <conflictor's-postinst> 'abort-remove' 'in-favour' <package>
#          <new-version>
#        * <postinst> 'abort-remove'
#        * <deconfigured's-postinst> 'abort-deconfigure' 'in-favour'
#          <failed-install-package> <version> 'removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package.

case "$1" in
    configure)
        # Create lirays user and group if they don't exist
        if ! getent group lirays > /dev/null 2>&1; then
            addgroup --system lirays --quiet
        fi

        if ! getent passwd lirays > /dev/null 2>&1; then
            adduser --system --quiet --group lirays --home /var/lib/lirays-scada
        fi

        # Create data directory
        mkdir -p /var/lib/lirays-scada/data_dir
        chown -R lirays:lirays /var/lib/lirays-scada
        chmod 755 /var/lib/lirays-scada
        chmod 755 /var/lib/lirays-scada/data_dir

        # Install systemd service
        systemctl daemon-reload

        # Enable and start the service
        systemctl enable lirays-scada.service
        systemctl start lirays-scada.service
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.



exit 0
